# edge-appliance/edge-stack/replication-edge.service
# =============================================================================
# systemd unit for the Replication Edge Appliance (Linux hosts & the VM image)
# Install:
#   sudo cp replication-edge.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now replication-edge
# =============================================================================
[Unit]
Description=Replication Edge Appliance (Netbird mesh + Kafka/Connect/Schema Registry)
Requires=docker.service netbird.service
After=docker.service netbird.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/replication-edge/edge-stack
# Brings up the mesh, resolves NETBIRD_IP, then `docker compose up -d`
ExecStart=/opt/replication-edge/edge-stack/start-edge.sh
# Graceful stop of the data plane (mesh is left to netbird.service)
ExecStop=/usr/bin/docker compose --env-file /opt/replication-edge/edge-stack/.env.edge down
TimeoutStartSec=300
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
